# cookbooks/boxcutter_prometheus/templates/nvidia_gpu_exporter.service.erb
[Unit]
Description=NVidia GPU Exporter
After=network-online.target
Wants=network-online.target

[Service]
User=prometheus
Group=prometheus
Type=simple
ExecStart=/opt/nvidia_gpu_exporter/latest/nvidia_gpu_exporter
ExecReload=/bin/kill -s SIGHUP $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target
