# cookbooks/boxcutter_prometheus/templates/pushgateway.service.erb
[Unit]
Description=Prometheus pushgateway
After=network-online.target
Wants=network-online.target

[Service]
User=prometheus
Group=prometheus
Type=simple
ExecStart=/opt/pushgateway/latest/pushgateway
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
