# cookbooks/boxcutter_prometheus/templates/blackbox_exporter/blackbox_exporter.service.erb
[Unit]
Description=Blackbox Exporter
Documentation=https://github.com/prometheus/blackbox_exporter
After=network-online.target
Wants=network-online.target

[Service]
User=prometheus
Group=prometheus
Type=simple
ExecReload=/bin/kill -s SIGHUP $MAINPID
ExecStart=/opt/blackbox_exporter/latest/blackbox_exporter \
<% options = node['boxcutter_prometheus']['blackbox_exporter']['command_line_flags'].reject { |k, _| k == 'config.file' } %>
<% options.each do |key, value| %>
  --<%= key %>="<%= value %>" \
<% end %>
  --config.file=/etc/blackbox_exporter/blackbox.yml
Restart=always

[Install]
WantedBy=multi-user.target
