# cookbooks/boxcutter_prometheus/templates/pushgateway/pushgateway.service.erb
[Unit]
Description=Prometheus pushgateway
After=network-online.target
Wants=network-online.target

[Service]
User=prometheus
Group=prometheus
Type=simple
<% flags = [] %>
<% node['boxcutter_prometheus']['redis_exporter']['command_line_flags'].each do |key, value| %>
  <%   flag = "--#{key}" %>
  <%   flag << "=#{value}" if value %>
  <%   flags << flag %>
<% end %>
ExecStart=/opt/pushgateway/latest/pushgateway <%= flags.join(' ') %>
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
